# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7.2]
command=/usr/sbin/php-fpm7.2
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
